# log2syslog.upstart
# log2syslog - generic log file feeder to a syslog server
#
# log2syslog reads one or more log files and send each log line to a remote
# syslog server, using UDP or TCP protocol.

description     "log2syslog, a generic log feeder"
author          "Klaubert Herr <klaubert@gmail.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

#setuid myuser
exec /usr/local/sbin/log2syslog -f /etc/log2syslog.conf
